# security_policies/36_cron_usage_audit/files/cron.rules
-w /usr/bin/crontab -p x -k cron_mod
-w /etc/cron.allow -p wa -k cron_allow
-w /etc/cron.deny -p wa -k cron_deny